# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:webserver]
command=nginx -g "daemon off;"
